# com.eventrelay.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.eventrelay</string>
    <key>ProgramArguments</key>
    <array>
        <string>/usr/local/bin/eventrelay</string>
        <string>--port</string>
        <string>6060</string>
    </array>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>StandardErrorPath</key>
    <string>/tmp/eventrelay.log</string>
    <key>StandardOutPath</key>
    <string>/tmp/eventrelay.log</string>
</dict>
</plist>
